# deploy/systemd/tecsol-parceiros.service
[Unit]
Description=Tecsol Parceiros (tp-app — TanStack Start com SSR)
After=network.target

[Service]
Type=simple
User=tecsol
WorkingDirectory=/opt/tecsol/tp-app
Environment=PORT=3000
Environment=NODE_ENV=production
ExecStart=/usr/bin/node .output/server/index.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
