# deploy/systemd/tecsol-api.service
[Unit]
Description=Tecsol API (Flask + Socket.IO)
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=tecsol
WorkingDirectory=/opt/tecsol/egn-tecsol/api
EnvironmentFile=/opt/tecsol/egn-tecsol/api/.env
# 1 worker só: o Socket.IO guarda as salas em memória; com vários workers o
# chat/presença quebraria sem um message queue compartilhado.
ExecStart=/opt/tecsol/venv/bin/gunicorn -w 1 --threads 8 -b 127.0.0.1:5000 run:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
